# roles/locksmithd/files/mockout.service
[Service]
ExecStop=/usr/bin/true
RemainAfterExit=yes
